# supervisor/x11vnc-xvfb.conf
[program:x11vnc-xvfb]
command=x11vnc -forever -usepw -create
user=root
startretries=3
autorestart=true
redirect_stderr=true
autostart=true
stopasgroup=true
killasgroup=true
